# old/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#

mount 192.168.7.177:/mnt/Starm /mnt/Starm
mount 192.168.7.102:/home/hyatt/Stor /mnt/Stor
