# QuadMas/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#

echo 1 > /proc/sys/kernel/sysrq
setkeycodes 75 43 # Model F '\ |' key mod
